# .config/gabriel_custom/backlight.rules
RUN+="/bin/chgrp wheel /sys/class/backlight/amdgpu_bl1/brightness"
RUN+="/bin/chmod g+w /sys/class/backlight/amdgpu_bl1/brightness"
